# macos/Users/example/Library/LaunchAgents/com.gitlab.stephen-fox.wgu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gitlab.stephen-fox.wgu</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/REPLACE-ME/go/bin/wgu</string>
    <string>up</string>
  </array>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE-ME/.wgu/wgu.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
